# install/kali-container/kali.container
# Kali Linux container (lolterm)
# Podman quadlet — installed to ~/.config/containers/systemd/kali.container
# Generates a systemd --user service that starts the container at boot.
#
# The container is privileged to support raw sockets, packet injection,
# and monitor-mode tools. Rootless podman --privileged is scoped to the
# user namespace and does not grant host-level root access.
#
# Selinux :Z on volume mounts is harmless in Permissive/Disabled mode
# and required in Enforcing mode for container read/write access.

[Unit]
Description=Kali Linux container (lolterm)
Documentation=https://github.com/r3b1s/lolterm
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/lolterm-kali
ContainerName=kali
Network=host
Privileged=true
Volume=/tmp/.X11-unix:/tmp/.X11-unix:Z
Volume=%h:%h:Z
Exec=sleep infinity

[Service]
Restart=always

[Install]
WantedBy=default.target
