# install/kali-container/kali.container
# Kali Linux container (lolterm)
# Podman quadlet — installed to ~/.config/containers/systemd/kali.container
# Generates a systemd --user service that starts the container at boot.
#
# The container itself is not privileged. Tool wrappers for tools that need
# raw sockets or packet injection pass --privileged at exec time, providing
# per-process escalation without elevating the entire container.
#
# SELinux separation is disabled for this container via SecurityLabelDisable=true.
# This allows the container to read/write mounted volumes (home directory, X11 socket)
# without SELinux relabeling. Rootless podman cannot relabel user_home_t to
# container_file_t on SELinux-enforcing systems, making :Z ineffective here.

[Unit]
Description=Kali Linux container (lolterm)
Documentation=https://github.com/r3b1s/lolterm
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/lolterm-kali
ContainerName=kali
Network=host
SecurityLabelDisable=true
Volume=/tmp/.X11-unix:/tmp/.X11-unix
Volume=%h:%h
Exec=sleep infinity

[Service]
Restart=always

[Install]
WantedBy=default.target
