# install/kali-container/kali.container
# Kali Linux container (lolterm)
# Podman quadlet — installed to ~/.config/containers/systemd/kali.container
# Generates a systemd --user service that starts the container at boot.
#
# The container itself is not privileged. Tool wrappers for tools that need
# raw sockets or packet injection pass --privileged at exec time, providing
# per-process escalation without elevating the entire container.
#
# SELinux :Z on volume mounts is harmless in Permissive/Disabled mode
# and required in Enforcing mode for container read/write access.

[Unit]
Description=Kali Linux container (lolterm)
Documentation=https://github.com/r3b1s/lolterm
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/lolterm-kali
ContainerName=kali
Network=host
Volume=/tmp/.X11-unix:/tmp/.X11-unix:Z
Volume=%h:%h:Z
Exec=sleep infinity

[Service]
Restart=always

[Install]
WantedBy=default.target
